# minecraft-server/crontab
* * * * * /home/ec2-user/check_connections.sh